# config/gemmy-monitor.conf
description "Gemmy monitor"

start on runlevel [2]
stop on runlevel [016]

expect daemon
exec bluepill load /var/www/gemmy/config/gemmy-monitor.pill

respawn